# packages/luci/luci-app-custom-splash/uci-defaults/device-defaults
#!/bin/sh

uci batch <<-EOF
	set freifunk.custom_splash='settings'
	set freifunk.custom_splash.mode='disabled'
	set freifunk.custom_splash.messages='disabled'
		commit freifunk
EOF
