# packages/luci/luci-app-custom-splash/uci-defaults/device-defaults
#!/bin/sh

touch /etc/config/custom_splash
uci batch <<-EOF
	set custom_splash.main='settings'
	set custom_splash.main.mode='disabled'
	set custom_splash.ain.messages='disabled'
		commit custom_splash
EOF
