# packages/utils/heartbeat/files/etc/uci-defaults/initialize_heartbeat
#!/bin/sh

#add heartbeat configuration
touch /etc/config/manager
uci get manager.heartbeat || uci set manager.heartbeat='heartbeat'
uci get manager.server || uci set manager.server='server'
uci get manager.heartbeat.enabled || uci set manager.heartbeat.enabled='1'
uci get manager.heartbeat.server || uci set manager.server.url='http://heartbeat.basicinside.de'

#function to print random integers
random () {
echo $((`dd if=/dev/urandom bs=4 count=1 2>/dev/null | hexdump -x | awk '$2 {print "0x" $2$3}'` % $1))
}

#add daily cron entry with random time
crontab -l|grep -q heart||echo `random 59` `random 23` \* \* \* /sbin/heartbeat>>/etc/crontabs/root
