# packages/utils/mac-hash/files/etc/init.d/machash
#!/bin/sh /etc/rc.common

START=55

start() {
    config_load manager
    config_get interface "heartbeat" "interface"
    for i in $interface ; do
	ifname=$(uci -p /var/state get network.$i.ifname)
	if [ -z "$OPTS" ] ; then
		OPTS="$ifname"
	else
		OPTS="$OPTS $ifname"
	fi
    done
    echo "Hearbeat machash listen on $OPTS"
    start-stop-daemon -S -x /sbin/machash -- "$OPTS"
}

stop() {
    start-stop-daemon -K -p /var/run/machash.pid
}

reload() { 
    start-stop-daemon -K -s 1 -p /var/run/machash.pid
} 
